# deploy/vps/calories-frontend.service
[Unit]
Description=Calories Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Calories/calories-frontend
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3001
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
